# systemd/.config/systemd/user/sync-repos.service
[Unit]
Description=Sync dotfiles and knowledge base repositories
After=network-online.target

[Service]
Type=oneshot
ExecStart=%h/bin/sync-repos
Environment="DOTFILES_DIR=%h/.dotfiles"
Environment="KNOWLEDGE_DIR=%h/brain"
